# deploy/critic-runner.service
[Unit]
Description=Nexus Critic Runner - Agent 自我表达引擎
After=network.target nexus.service
Wants=nexus.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/nexus
Environment=PYTHONPATH=/opt/nexus
ExecStart=/usr/bin/python3 -m src.critic_runner
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
